# agent/scripts/wira.service
[Unit]
Description=Wira — AI assistant for WhatsApp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=craig
Group=craig
WorkingDirectory=/home/craig/claude-creations/wira-whatsapp
EnvironmentFile=/home/craig/claude-creations/wira-whatsapp/.env
ExecStart=/home/craig/claude-creations/wira-whatsapp/.venv/bin/python main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
