# agent/scripts/wira.service
[Unit]
Description=Wira — AI assistant for WhatsApp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=__AGENT_DIR__
EnvironmentFile=__AGENT_DIR__/.env
ExecStart=__AGENT_DIR__/.venv/bin/python main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
